# bbs.conf
[program:bbs]
command=/usr/local/bin/gunicorn wsgi --bind 0.0.0.0:2001 --pid /tmp/bbs.pid
directory=/var/www/bbs
autostart=true
autorestart=true
