# bbs.conf
[program:bbs]
command=/usr/local/bin/gunicorn wsgi -c gunicorn_conf.py
directory=/var/www/bbs
autostart=true
autorestart=true
